# srcs/requirements/wordpress/conf/php-fpm.conf
[supervisord]
nodaemon=true
[program:php-fpm]
command=/usr/sbin/php-fpm7.3 --nodaemonize
autostart=true
autorestart=true
